# standby-monitor.service
[Unit]
Description=Run script after wake or boot

[Service]
ExecStart={{SCRIPT_PATH}}

[Install]
WantedBy=sleep.target multi-user.target