# standby-monitor.service
[Unit]
Description=Run script after wake

[Service]
ExecStart=/path/to/your/script.sh

[Install]
WantedBy=sleep.target